# foot-server.socket
[Socket]
ListenStream=%t/foot-wayland-1.sock
Service=foot-server.service

[Install]
WantedBy=sockets.target
